# uwrm.yaml
title: Unauthorized Windows Registry Modification
id: 78901234-56ab-78cd-90ef-1234567890ab
description: Detects modifications to sensitive registry keys that could indicate persistence or privilege escalation.
status: experimental
author: -pk
date: 2025/05/16
logsource:
  category: registry_event
  product: windows
detection:
  selection:
    EventType: 'SetValue'
    TargetObject|contains:
      - '\CurrentVersion\Run'
      - '\CurrentVersion\RunOnce'
      - '\CurrentVersion\Image File Execution Options'
  condition: selection
fields:
  - TargetObject
  - Details
falsepositives:
  - Legitimate software installations (e.g., Windows Updates, third-party installers)
level: medium
tags:
  - attack.persistence
  - attack.t1547.001
